# web/cron/crontab
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

* * * * * root /opt/haera/organize.sh >> /var/log/haera.log 2>&1
